# ansible/roles/monitoring/templates/node_exporter.service.j2
# {{ ansible_managed }}
[Unit]
Description=Prometheus Node Exporter
Documentation=https://github.com/prometheus/node_exporter
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User={{ node_exporter_user }}
Group={{ node_exporter_group }}
ExecStart={{ node_exporter_install_dir }}/node_exporter \
  --web.listen-address=:{{ node_exporter_port }} \
  {{ node_exporter_additional_args }}

Restart=on-failure
RestartSec=5s

[Install]
WantedBy=multi-user.target
